# files/ojs_daily.cron
#!/usr/bin/env bash

PATH=/opt/oulib/ojs/bin:/usr/local/bin:/usr/bin:/bin:/sbin:$PATH

printf "Searching for OJS sites.\n"
find -L /srv -maxdepth 3 -mindepth 3 -wholename "*/sites/*" -type d 2>/dev/null | while read -r SITEPATH ; do
  # Check for config.inc.php
  stat ${SITEPATH}/ojs/config.inc.php &>/dev/null
  if [ "$?" -eq "0" ]; then
    printf "Backing up ${SITEPATH}\n"
    /opt/oulib/ojs/bin/ojs_snapshot.sh ${SITEPATH}
  fi
done

printf "Done.\n"
